# remotemanager_server/system/rootfs/buildroot-2014.02/package/ntp/ntpdate.service
[Unit]
Description=Update system time with NTP

[Service]
Type=simple
ExecStart=/usr/sbin/ntpd -g -n -q
