# remotemanager_server/system/rootfs/buildroot-2014.02/package/remotemanager/rfcommwatchdog.service
[Unit]
Description=rfcomm watchdog

[Service]
ExecStart=/usr/bin/python  /srv/remotemanager/rfcomm_watchdog.py /data/rfcomm-watchdog.log

[Install]
WantedBy=multi-user.target
